# etc/init/nas.conf
# nas - Bring up the NAS

description	"NAS script"

start on (startup
          and started udev
          and stopped udevtrigger)

exec /bin/nas start
